# wallpaper_shuffle/wallpaper-shuffle.timer
[Unit]
Description=Shuffle wallpaper every hour
After=graphical-session.target

[Timer]
OnBootSec=30sec
OnUnitActiveSec=1h
Persistent=true
AccuracySec=1s

[Install]
WantedBy=default.target
